# container/consultainer.conf
[program:webserver]
command=/usr/sbin/apache2ctl -DFOREGROUND
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true

